# autobackup.service
[Unit]
Description=Autobackup

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/sbin/backup.sh
WorkingDirectory=/tmp

